# templates/plugins/telegram/telegram.supervisord.conf
########################################################################
# telegram
########################################################################
[group:telegram]
programs=cli
#programs=cli,bot,intouch

[program:cli]
directory={{ redmine_base_dir }}
command=/usr/local/bin/telegram-cli
    --username redmine
    --disable-readline
    --disable-output
    --disable-colors
    --wait-dialog-list
    --json
    --tcp-port 2391
    --logname {{ redmine_base_dir }}/log/telegram-cli.log
user=root
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=syslog

#[program:bot]
#directory={{ redmine_base_dir }}
#environment=HOME="{{ redmine_user_home_dir }}",PID_DIR="/tmp",RAILS_ENV="production"
#command=/usr/local/bin/bundle exec
#    rake
#    chat_telegram:bot
#user={{ redmine_user }}
#autostart=true
#autorestart=true
#redirect_stderr=True
#stdout_logfile=syslog

#[program:intouch]
#directory={{ redmine_base_dir }}
#environment=HOME="{{ redmine_user_home_dir }}",PID_DIR="/tmp",RAILS_ENV="production"
#command=/usr/local/bin/bundle exec
#    rake
#    intouch:telegram:bot
#user={{ redmine_user }}
#autostart=true
#autorestart=true
#redirect_stderr=True
#stdout_logfile=syslog
